# ops/ec2/thathinh.service
[Unit]
Description=Tha Thinh API (Spring Boot)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/opt/thathinh
EnvironmentFile=/opt/thathinh/env.sh
ExecStart=/usr/bin/java -jar /opt/thathinh/app.jar --server.port=8081 --spring.profiles.active=prod
Restart=on-failure
RestartSec=10
SuccessExitStatus=143
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
